# plcwatchd.draft
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

home=/opt/plcwatchd
command="${home}/plcwatchd"
command_args="-i 192.168.178.105 -k xxxxxxxxxxxxxxxx -t xxxxxxxxxxxxxxxxx -c 600 -e 3600 -l
/var/log/plcwatchd.log"
pidfile=/run/plcwatchd.pid
start_stop_daemon_args="--chdir ${home} -p ${pidfile} --background -m"

depend() {
  need net
}

